# build/cron/crontab
*/15 * * * * php bin/console schedule:check:date